# installation_files/covid-dashboard-filewatcher.path
[Path]
PathModified=/home/ubuntu/covid19-dashboard/data/COVID-19_aantallen_gemeente_cumulatief.csv

[Install]
WantedBy=multi-user.target

